# backend/deploy/mundografic@.service
# systemd template unit for Mundografic (blue/green instances)
# Save as /etc/systemd/system/mundografic@.service
#
# Instancias disponibles:
#   mundografic@green.service  → copia staging
#
# El .env de cada instancia se carga desde el WorkingDirectory/backend/.env

[Unit]
Description=MundoGrafic backend — instancia %i
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mundografic/%i/backend

# NO usar EnvironmentFile — dotenv carga el .env desde WorkingDirectory automáticamente.
# Esto evita problemas con caracteres especiales (@, #, etc.) en el parser de systemd.

# Ejecuta el servidor compilado (npm run build debe haberse ejecutado antes)
ExecStart=/usr/bin/node /opt/mundografic/%i/backend/dist/server.js

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Límites de recursos (ajustable según tu servidor)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
